# supervisord.conf
[supervisord]

[program:haldibot]
directory=/home/ec2-user/haldibot./
command=python3 haldibot.py
autostart=true
autorestart=true